# package/linux/debian-control/postinst.in
#!/bin/sh

# errors shouldn't cause script to exit
set +e 

# add rserver user account
sudo useradd -r rstudio-server

# create softlink to admin script in /usr/sbin
sudo ln -f -s ${CMAKE_INSTALL_PREFIX}/bin/rstudio-server /usr/sbin/rstudio-server

# create config directory
sudo mkdir -p /etc/rstudio

# create var directories
sudo mkdir -p /var/run/rstudio-server
sudo mkdir -p /var/lock/rstudio-server
sudo mkdir -p /var/log/rstudio-server
sudo mkdir -p /var/lib/rstudio-server
sudo mkdir -p /var/lib/rstudio-server/conf
sudo mkdir -p /var/lib/rstudio-server/body
sudo mkdir -p /var/lib/rstudio-server/proxy

# write installed indicator then force suspend active sessions
# (so that the new client is never paired against the old server)
sudo mkdir -p /etc/rstudio
sudo sh -c "echo `date +%s` > /var/lib/rstudio-server/installed"
sudo rstudio-server force-suspend-all

# check lsb release
LSB_RELEASE=`lsb_release --id --short`

# add apparmor profile
if test $LSB_RELEASE = "Ubuntu" && test -d /etc/apparmor.d/
then
   sudo cp ${CMAKE_INSTALL_PREFIX}/extras/apparmor/rstudio-server /etc/apparmor.d/
fi

# add upstart profile or init.d script and start the server
if test $LSB_RELEASE = "Ubuntu" && test -d /etc/init/
then
   sudo cp ${CMAKE_INSTALL_PREFIX}/extras/upstart/rstudio-server.conf /etc/init/
   sudo initctl reload-configuration
   sudo initctl stop rstudio-server 2>/dev/null
   sudo initctl start rstudio-server
else
   sudo cp ${CMAKE_INSTALL_PREFIX}/extras/init.d/debian/rstudio-server /etc/init.d/
   sudo chmod +x /etc/init.d/rstudio-server
   sudo update-rc.d rstudio-server defaults
   sudo /etc/init.d/rstudio-server stop  2>/dev/null
   sudo /etc/init.d/rstudio-server start
fi

# clear error termination state
set -e
