# src/cpp/server/extras/systemd/rstudio-server.service.in
[Unit]
Description=RStudio Server

[Service]
Type=forking
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/rserver
ExecStop=/usr/bin/killall -TERM rserver
KillMode=none
Restart=on-failure

[Install]
WantedBy=multi-user.target
