# src/cpp/server/extras/systemd/rstudio-server.redhat.service.in
[Unit]
Description=RStudio Server

[Service]
Type=forking
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/rserver
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
