# src/cpp/server/extras/systemd/rstudio-server.service.in
[Unit]
Description=RStudio Server
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/run/rstudio-server.pid
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/rserver
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
